# supervisord.conf
[program:scrapy_daemon]
command = /home/dotcloud/current/scrapyd.sh
